# JCatascopia-Agent/scripts/JCatascopia-Agent
### BEGIN INIT INFO
# Provides:          JCatascopia-Server
# Required-Start:    
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: JCatascopia-Agent deamon
# Description:       JCatascopia Monitoring Agent v0.2 init.d script
### END INIT INFO

# Author: Demetris Trihinas <trihinas@cs.ucy.ac.cy>

NAME=JCatascopia-Agent
LOCKFILE=/var/lock/JCatascopia-Agent-lock
JCPATH=/usr/local/bin/JCatascopiaAgentDir

case $1 in
    start)
	if [ -e $LOCKFILE ] ; then
            echo "$NAME already up and running..."
    	    exit
        fi
        touch $LOCKFILE
        /bin/bash $JCPATH/JCatascopia-Agent-start.sh
        echo "Starting $NAME..."
    ;;
    stop)
        if [ -e $LOCKFILE ] ; then
            rm -f $LOCKFILE
            /bin/bash $JCPATH/JCatascopia-Agent-stop.sh
            echo "Stopping $NAME..."
            exit
        fi
	    echo "$NAME not running..."
    ;;
    restart)
	$0 stop
	$0 start
    ;;
    *)
    	echo "Usage: /etc/init.d/$NAME {start|stop|restart}"
esac
exit
